# etc/supervisor/conf.d/crashplan.conf
[program:crashplan]
command=/etc/init.d/crashplan start
user=root
autorestart=true
priority=400